# supervisord/api.conf
# supervisord -c supervisord/api.conf
[supervisord]
nodaemon = true

[program:backend]
command=./build/run backend:devserver
numprocs=1
autostart=true
redirect_stderr=true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0

[program:indexer]
command=./build/run indexer
numprocs=1
autostart=true
redirect_stderr=true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
