# supervisord/backend.conf
# supervisord -c supervisord/backend.conf
[supervisord]
nodaemon = true

[program:backend]
command=./build/run backend:devserver
numprocs=1
autostart=true
redirect_stderr=true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
